# openshift-offline-resources/okd-v3.11.0-0cbc58b/resource/openshift-ansible/roles/openshift_gcp/files/openshift-bootstrap-update.service
[Unit]
Description=Update the OpenShift node bootstrap configuration

[Service]
Type=oneshot
ExecStart=/usr/bin/openshift-bootstrap-update
User=root
