# api/resource/packages/okd-v3-11-0/openshift-ansible/roles/openshift_gcp/files/openshift-bootstrap-update.timer
[Unit]
Description=Update the OpenShift node bootstrap credentials hourly

[Timer]
OnBootSec=30s
OnCalendar=hourly
Persistent=true     
 
[Install]
WantedBy=timers.target